# rtaudio/target/rtaudio-wifi.service
[Unit]
Description=RT-Audio WiFi Service
After=network.target rtaudio-config.service

[Service]
Type=simple
ExecStart=/home/rtaudio/rtaudio-wifi.sh
KillSignal=SIGINT
Restart=on-failure
RestartSec=4

[Install]
WantedBy=multi-user.target
